# apps/haproxy-supervisor/docker/confd-supervisor.conf
[program:confd]
command=/usr/local/sbin/confd -f -interval=10 -node=127.0.0.1:2379  /usr/local/bin/etc/haproxy/haproxy.cfg
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
