# apps/haproxy-supervisor/docker/haproxy-supervisor.ini
[program:haproxy]
command=/usr/local/sbin/haproxy -f /etc/haproxy/haproxy.cfg
autostart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
